# packages/jaraco.home/jaraco_home-3.10.2.tar.gz/jaraco_home-3.10.2/jaraco/home/Gather HDHomeRun Stats.plist
<?xml version="1.0" encoding="UTF-8"?>

<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">

<plist version="1.0">

<dict>
	<key>Label</key>
	<string>com.jaraco.home.hdhomerun</string>
	<key>ProgramArguments</key>
	<array>
		<string>{sys.executable}</string>
		<string>-m</string>
		<string>jaraco.home.hdhomerun</string>
		<string>--update</string>
	</array>
	<key>StartInterval</key>
	<integer>3600</integer>
	<key>SessionCreate</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>{logs}/hdhomerun.err</string>
	<key>StandardOutPath</key>
	<string>{logs}/hdhomerun.out</string>
</dict>

</plist>
